# serveo.service
[Unit]
Description=SSH tunneling

[Service]
Type=simple
ExecStart=/bin/bash HOME/serveo/start.sh

[Install]
WantedBy=multi-user.target
